# configs/unattended/root-image/etc/systemd/system/multi-user.target.wants/arch-install.service
[Unit]
Description=Arch Install

[Service]
ExecStart=/usr/local/bin/arch-install

[Install]
WantedBy=multi-user.target
