# configs/baseline/root-image/etc/rc.local
#!/bin/sh

set -e

printf 'o\nY\nn\n\n\n\n\nw\nY\n' | gdisk "/dev/sda"
mkfs -t"ext4" "/dev/sda1"
mount "/dev/sda1" "/mnt"
pacstrap "/mnt" "base"
genfstab -p "/mnt" >>"/mnt/etc/fstab"
arch-chroot "/mnt" "/usr/local/bin/arch-chroot-install"
umount "/mnt"
reboot
